# debian/extra/units/systemd-fsck-root.service.d/10-skip-fsck-initramfs.conf
[Unit]
ConditionPathExists=!/run/initramfs/fsck-root
